# files/etc/crontabs/root
0 * * * * sh /root/schedule/restart-wireguard.sh >> /var/log/restart-wireguard.log 2>&1
