# files/etc/crontabs/root
* * * * * sh /root/schedule/facebook-re.sh >> /var/log/facebook-re.log 2>&1
* * * * * sh /root/schedule/google-re.sh >> /var/log/google-re.log 2>&1
* * * * * sh /root/schedule/twitter-re.sh >> /var/log/twitter-re.log 2>&1
10 3 * * 5 sh /root/schedule/adlist-update.sh >> /var/log/adlist-update.log 2>&1
0 3 * * 5 lua /usr/share/passwall/rule_update.lua log > /dev/null 2>&1 &