# files/etc/crontabs/root
*/15 * * * * sh /root/schedule/facebook-re.sh >> /var/log/facebook-re.log 2>&1
*/15 * * * * sh /root/schedule/google-re.sh >> /var/log/google-re.log 2>&1
*/15 * * * * sh /root/schedule/twitter-re.sh >> /var/log/twitter-re.log 2>&1
0 * * * * sh /root/schedule/restart-wireguard.sh >> /var/log/restart-wireguard.log 2>&1