# files/etc/crontabs/root
0 * * * * sh /root/schedule/restart-wireguard.sh >> /var/log/restart-wireguard.log 2>&1
10 3 * * 5 sh /root/schedule/adlist-update.sh >> /var/log/adlist-update.log 2>&1